# supervisord.conf
[supervisord]
nodaemon=true

[program:emulator]
command=python /app/emulator.py

[program:server]
command=node /app/server.js